# assignments/sqlite-mcp-gui/templates/bootstrap/hetzner-codespaces.service
[Unit]
Description=Genesis Application Server (com.hetzner.codespaces)
Documentation=https://github.com/ebowwa/com.hetzner.codespaces
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=genesis
Group=genesis
WorkingDirectory=/opt/genesis

# Execution
ExecStart=/usr/bin/bun start
ExecReload=/bin/kill -HUP $MAINPID

# Restart Policy (with rate limiting)
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=300   # 5-minute window
StartLimitBurst=5           # Max 5 restarts

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=genesis

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=/etc/default/genesis
EnvironmentFile=-/etc/default/genesis.local  # Optional override (- = not required)

# Resource Limits
LimitNOFILE=65536

# Security Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/genesis/data /var/log/genesis

[Install]
WantedBy=multi-user.target
